# systemd/apps/letsencrypt_renew.service
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=/root/certbot-auto renew --quiet --agree-tos
